# supervisord.conf
# supervisord.ini file for mediacloud web-search this launches and
# supervises all process_tasks workers in a single container.
[supervisord]
nodaemon=true

[program:system-slow-worker]
command= python mcweb/manage.py process_tasks --queue system-slow

[program:system-fast-worker]
command= python mcweb/manage.py process_tasks --queue system-fast

[program:admin-slow-worker]
command= python mcweb/manage.py process_tasks --queue admin-slow

[program:admin-fast-worker]
command= python mcweb/manage.py process_tasks --queue admin-fast

[program:user-slow-worker]
command= python mcweb/manage.py process_tasks --queue user-slow

[program:user-fast-worker]
command= python mcweb/manage.py process_tasks --queue user-fast
