# periodic/daily/300.task-log-cleanup
#!/bin/sh

# cleanup task and old rescrape logs older than 30 days

for DIR in rescrape tasks; do
    PATH=data/logs/$DIR
    if [ -d $PATH ]; then
	find $PATH -type f -mtime +30 -delete
    fi
done
