# periodic/weekly/110.update-pub-date
#!/bin/sh

python mcweb/manage.py sources-meta-update --queue  --task publication_date