# periodic/weekly/120.source-alert-system
#!/bin/sh

python mcweb/manage.py source-alert-system --update --queue
