# periodic/weekly/250.delete-reset-tokens
#!/bin/sh

python mcweb/manage.py delete-reset-tokens