# periodic/weekly/400.update-stories-per-week
#!/bin/sh

python mcweb/manage.py sources-meta-update --queue --update --task stories_per_week
